# orangepi-zero/overlayfs/etc/profile.d/welcomemsg.sh
cat /etc/profile.d/message.txt

neofetch
